# ffulm-migration/files/lib/gluon/ffulm-migration/rc.d/S13led
#!/bin/sh /etc/rc.common

# shellcheck disable=SC2034
START=13

start() {
        /etc/init.d/led start
        # shellcheck source=/dev/null
        . /etc/diag.sh
        set_state failsafe
}
